# packaging/openscanstation-watchdog.timer
[Unit]
Description=OpenScanStation regelmäßig prüfen

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Persistent=true
Unit=openscanstation-watchdog.service

[Install]
WantedBy=timers.target
