# packaging/openscanstation-device-settings.service
[Unit]
Description=OpenScanStation Geräteeinstellungen
After=network.target openscanstation.service
Wants=openscanstation.service

[Service]
Type=simple
User=openscanstation
Group=openscanstation
ExecStart=/usr/bin/openscanstation-device-settings --host 0.0.0.0 --port 8102
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/openscanstation

[Install]
WantedBy=multi-user.target
