# packaging/openscanstation-watchdog.service
[Unit]
Description=OpenScanStation Health-Check
After=openscanstation.service

[Service]
Type=oneshot
ExecStart=/usr/bin/openscanstation-watchdog
